# config/includes.chroot/etc/profile.d/zurvan.sh
# Zurvan Linux environment branding
export DISTRIB_ID=Zurvan
export DISTRIB_RELEASE=0.1.0
export DISTRIB_CODENAME=akaran
export DISTRIB_DESCRIPTION="Zurvan Linux 0.1.0"
