# config/includes.chroot/etc/systemd/system/spice-vdagentd.service.d/override.conf
[Unit]
Description=Spice vdagent

[Service]
Type=forking
ExecStart=/usr/sbin/spice-vdagent -d

[Install]
WantedBy=multi-user.target
