# config/includes.chroot/etc/systemd/system/open-vm-tools.service.d/override.conf
[Unit]
Description=Open Virtual Machine Tools Service

[Service]
Type=forking
ExecStart=/usr/sbin/open-vm-tools --daemon

[Install]
WantedBy=multi-user.target
