# config/includes.chroot/etc/systemd/system/qemu-guest-agent.service.d/override.conf
[Unit]
Description=QEMU Guest Agent

[Service]
Type=forking
ExecStart=/usr/sbin/qemu-ga --daemon

[Install]
WantedBy=multi-user.target
